# servers.d/leo.com.br.config
#!/bin/sh
#
# chkconfig:   - 85 15

RUBY_VERSION="2.6.0"
APP_DIR="/var/www/rinvest"

source /home/deploy/puma_servers/base
